# supporting_files/supervisord-cronjob.conf
[program:cronjob]
command=cron -f
autostart=true
autorestart=true
stderr_logfile=/var/log/cronjob.err.log
stdout_logfile=/var/log/cronjob.out.log